# com.ima2obsidian.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ima2obsidian.update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/berton/Github/Ima2Obsidian/ima_incremental_update.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/berton/Github/Ima2Obsidian/incremental_update.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/berton/Github/Ima2Obsidian/incremental_update_error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/berton/Github/Ima2Obsidian</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
